# scheduler/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# m h dom mon dow user  command
*/2 * * * * root /app/run.sh >> /var/log/cron.log 2>&1
